# supervisord.conf
[program:oauth2-proxy]
directory=/app
command=""
autostart=true
autorestart=true
redirect_stderr=true